# systemd/agnetic-nats.service
[Unit]
Description=Starship OS - NATS Message Bus
After=network.target
Documentation=https://github.com/andromi-hash/starship-os

[Service]
Type=simple
# active.conf → agent-bus.conf (dev/server/edge) or fleet-bus.conf (ops)
ExecStart=/usr/local/bin/nats-server -c /etc/starship/nats/active.conf
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
User=nats
Group=nats
EnvironmentFile=-/etc/starship/nats.env
LimitNOFILE=65535
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/starship/nats /var/lib/agnetic/nats /var/log/starship /var/log/agnetic
ProtectHome=true
KillMode=mixed
TimeoutStopSec=30
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
