# deploy/agnetic-dashboard.target
[Unit]
Description=Starship OS Dashboard — status bridge + web UI
After=nats.target
BindsTo=agnetic-agent-mesh.target
Wants=agnetic-status-bridge.service agnetic-dashboard-web.service

[Install]
WantedBy=agnetic-agent-mesh.target
