# deploy/nats.service
[Unit]
Description=NATS Server for Starship OS Agent Bus
After=network.target

[Service]
Type=simple
User=tech
ExecStart=/usr/sbin/nats-server -c /home/tech/starship-os/nats/agent-bus.conf
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
