# debian/DEBIAN/postinst
#!/bin/bash
# Starship OS postinst — after package extract
set -e

echo "╔══════════════════════════════════════════════╗"
echo "║  Starship OS — Post-Install Setup   ║"
echo "╚══════════════════════════════════════════════╝"
echo ""

# ─── 1. Create system users ─────────────────────────────────────────
if ! id -u agnetic &>/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin \
        --comment "Starship OS Service Account" agnetic
    echo "[OK] Created user: agnetic"
else
    echo "[OK] User agnetic already exists"
fi

if ! id -u nats &>/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin \
        --comment "NATS Message Bus" nats
    echo "[OK] Created user: nats"
else
    echo "[OK] User nats already exists"
fi

# ─── 2. Directories + legacy symlinks ────────────────────────────────
mkdir -p /opt/starship/bin /opt/starship/lib/starship /opt/starship/venv
mkdir -p /etc/starship/nats /var/lib/starship/nats /var/lib/starship/message-history /var/log/starship
ln -sfn /opt/starship /opt/agnetic
ln -sfn /etc/starship /etc/agnetic
ln -sfn /var/lib/starship /var/lib/agnetic
ln -sfn /var/log/starship /var/log/agnetic
# active.conf default if missing
if [[ ! -e /etc/starship/nats/active.conf ]]; then
    ln -sfn /etc/starship/nats/agent-bus.conf /etc/starship/nats/active.conf
fi
echo "[OK] Directories created"

# ─── 3. Python venv ─────────────────────────────────────────────────
if [[ ! -x /opt/starship/venv/bin/python3 ]]; then
    if python3 -m venv /opt/starship/venv 2>/dev/null; then
        echo "[OK] Created venv"
    else
        echo "[WARN] python3 -m venv failed (install python3-venv)"
    fi
fi
if [[ -x /opt/starship/venv/bin/pip ]]; then
    /opt/starship/venv/bin/pip install --upgrade pip -q 2>/dev/null || true
    /opt/starship/venv/bin/pip install nats-py aiohttp httpx PyYAML -q 2>/dev/null || true
    echo "[OK] Python venv ready"
else
    echo "[WARN] venv pip missing — agents need: python3 -m venv /opt/starship/venv && pip install nats-py aiohttp httpx PyYAML"
fi

# ─── 4. Permissions ─────────────────────────────────────────────────
chown -R agnetic:agnetic /opt/starship /etc/starship /var/lib/starship /var/log/starship 2>/dev/null || true
chown -R nats:nats /var/lib/starship/nats 2>/dev/null || true
# NATS conf readable by nats user
chown -R root:nats /etc/starship/nats 2>/dev/null || true
chmod 755 /etc/starship/nats 2>/dev/null || true
chmod +x /opt/starship/bin/* 2>/dev/null || true
echo "[OK] Permissions set"

# ─── 5. Systemd ─────────────────────────────────────────────────────
if command -v systemctl &>/dev/null; then
    systemctl daemon-reload || true
    for unit in agnetic-nats agnetic-staragent \
                agnetic-agent@proxy agnetic-agent@romi agnetic-agent@ergo \
                agnetic-status-bridge agnetic-message-history agnetic-dashboard \
                starship-fleet; do
        systemctl enable "${unit}.service" 2>/dev/null || true
    done
    systemctl enable agnetic-mesh.target 2>/dev/null || true
    echo "[OK] Services enabled"
fi

# ─── 6. GPU detection (optional) ────────────────────────────────────
if [[ -x /opt/starship/bin/detect-gpu.sh ]]; then
    /opt/starship/bin/detect-gpu.sh detect 2>/dev/null || true
fi

echo ""
echo "Installation complete!"
echo ""
echo "  Firstboot: STARSHIP_PROFILE=ops /opt/starship/bin/starship-firstboot.sh"
echo "  Start:     systemctl start agnetic-mesh.target"
echo "  CLI:       starshipctl fleet status"
echo "  Dashboard: http://localhost:8788"
echo ""
exit 0
